# fcgiwrap.conf
[program:fcgiwrap]
user=www-data
command=/usr/sbin/fcgiwrap -f -s unix:/tmp/fcgiwrap.socket
autorestart=true
