# supervisor-fcgiwrap.conf
[program:fcgiwrap]
user=nexus
command=/usr/sbin/fcgiwrap -f -s unix:/tmp/fcgiwrap.socket
autorestart=true
